# deploy/auto_coin.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
    auto_coin launchd 템플릿 (macOS).

    설치:
        1) 이 파일을 홈 LaunchAgents로 복사:
             cp deploy/auto_coin.plist ~/Library/LaunchAgents/
        2) WorkingDirectory / 사용자 홈 경로를 본인 환경에 맞게 수정
             (이 파일의 __PROJECT_DIR__ 과 __HOME__ 자리표시자를 실제 값으로 교체)
             예) sed -i '' "s|__PROJECT_DIR__|$HOME/IdeaProjects/auto_coin|g; s|__HOME__|$HOME|g" \
                     ~/Library/LaunchAgents/auto_coin.plist
        3) 부팅 시 자동 시작 + 즉시 로드:
             launchctl unload ~/Library/LaunchAgents/auto_coin.plist 2>/dev/null
             launchctl load   ~/Library/LaunchAgents/auto_coin.plist

    확인:
        launchctl list | grep auto_coin
        tail -f __PROJECT_DIR__/logs/launchd.out.log

    정지 / 제거:
        launchctl unload ~/Library/LaunchAgents/auto_coin.plist
-->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>auto_coin</string>

    <key>ProgramArguments</key>
    <array>
        <string>__PROJECT_DIR__/.venv/bin/python</string>
        <string>-m</string>
        <string>auto_coin.web</string>
        <string>--host</string>
        <string>0.0.0.0</string>
        <string>--port</string>
        <string>8080</string>
    </array>

    <key>WorkingDirectory</key>
    <string>__PROJECT_DIR__</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>HOME</key>
        <string>__HOME__</string>
        <key>PATH</key>
        <string>/usr/local/bin:/opt/homebrew/bin:/usr/bin:/bin</string>
        <key>LANG</key>
        <string>en_US.UTF-8</string>
    </dict>

    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>

    <key>StandardOutPath</key>
    <string>__PROJECT_DIR__/logs/launchd.out.log</string>
    <key>StandardErrorPath</key>
    <string>__PROJECT_DIR__/logs/launchd.err.log</string>

    <!-- 실패 시 재시작 간격 (초) -->
    <key>ThrottleInterval</key>
    <integer>10</integer>
</dict>
</plist>
